# misc/99-hdmi.rules
ACTION=="change", SUBSYSTEM=="myHDMI", RUN+="/bin/sh /home/pi/send_power_ir.sh"